# systemd/services/web-deployd.timer
[Unit]
Description=Auto-deploy website every 5 seconds

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s

[Install]
WantedBy=timers.target
