# systemd/services/web-deployd.service
[Unit]
Description=Website deploy service

[Service]
ExecStart=/usr/local/bin/deploy-web
